# sample.crontab
LOGDIR=/tmp

@reboot	sleep 10; sudo pigpiod
@reboot	sleep 11; $HOME/bin/MusicBox wsserver > $LOGDIR/wsserver.log 2>&1
@reboot	sleep 12; $HOME/bin/MusicBox wsserver -w 1 -p 8881> $LOGDIR/wsserver.log 2>&1
@reboot	sleep 13; $HOME/bin/MusicBox wsserver -w 2 -p 8882 > $LOGDIR/wsserver.log 2>&1
@reboot	sleep 14; $HOME/bin/MusicBox wsserver -w 3 -p 8883 > $LOGDIR/wsserver.log 2>&1
@reboot	sleep 15; $HOME/bin/MusicBox calibration > $LOGDIR/calibration.log 2>&1